# packaging/linux/beetle.service
[Unit]
Description=Beetle agent (甲壳虫)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Create a dedicated user on the device and set User=/Group= here (see README.txt).
# User=beetle
# Group=beetle
WorkingDirectory=/opt/beetle/current
ExecStart=/opt/beetle/current/beetle
Restart=on-failure
RestartSec=5
# State and SPIFFS-emulation paths (host): see README.txt for BEETLE_STATE_ROOT.
ReadWritePaths=/var/lib/beetle /data/beetle
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
